# update_dns/update_dns.service
[Unit]
Description=Update GoDaddy DNS Record

[Service]
EnvironmentFile=/etc/systemd/system/update_dns.env
ExecStart=/usr/local/bin/update_dns.sh
