# update_dns/update_dns.timer
[Unit]
Description=Periodically update DNS with local IPv6 address

[Timer]
OnBootSec=30s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
