# update_dns/update_dns.service
[Unit]
Description=Update namecheap DNS Record

[Service]
EnvironmentFile=/etc/systemd/update_dns.env
ExecStart=/usr/local/bin/update_dns.sh
